# scripts/crontab
* * * * * /opt/bookcars/scripts/noip.sh 2>&1 | /usr/bin/logger -t noip-cron
@reboot /opt/bookcars/scripts/noip.sh 2>&1 | /usr/bin/logger -t noip-cron
